# shadowsocks/files/ssmgr.service
[Unit]
Description=Daemon to start ssmgr
Wants=network-online.target
After=network.target

[Service]
User=shadowsocks
Group=shadowsocks
Type=simple
ExecStart=/usr/bin/ssmgr -c ss.yml
Restart=on-abort

[Install]
WantedBy=multi-user.target
